# deploy/vietnix/jdm-queue.service
[Unit]
Description=JDM WORLD Laravel queue worker
After=network-online.target mysql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/jdm-world/backend
ExecStart=/usr/bin/php artisan queue:work --sleep=3 --tries=3 --timeout=120 --max-time=3600
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
